# awsvpnclient/reference.file.list
/opt/awsvpnclient
/opt/awsvpnclient/ACVC.GTK.Service
/opt/awsvpnclient/ACVC.GTK.Service.deps.json
/opt/awsvpnclient/ACVC.GTK.Service.dll
/opt/awsvpnclient/ACVC.GTK.Service.runtimeconfig.json
/opt/awsvpnclient/AWS VPN Client.deps.json
/opt/awsvpnclient/AWS VPN Client.dll
/opt/awsvpnclient/AWS VPN Client.runtimeconfig.json
/opt/awsvpnclient/AWSSDK.CloudWatch.dll
/opt/awsvpnclient/AWSSDK.CognitoIdentity.dll
/opt/awsvpnclient/AWSSDK.Core.dll
/opt/awsvpnclient/AWSSDK.Pinpoint.dll
/opt/awsvpnclient/AWSSDK.S3.dll
/opt/awsvpnclient/AWSSDK.SecurityToken.dll
/opt/awsvpnclient/AWSVPNClient
/opt/awsvpnclient/AWSVPNClient.Core.dll
/opt/awsvpnclient/AtkSharp.dll
/opt/awsvpnclient/Autofac.dll
/opt/awsvpnclient/CairoSharp.dll
/opt/awsvpnclient/EntityFramework.SqlServer.dll
/opt/awsvpnclient/EntityFramework.dll
/opt/awsvpnclient/GLibSharp.dll
/opt/awsvpnclient/GdkSharp.dll
/opt/awsvpnclient/GioSharp.dll
/opt/awsvpnclient/GtkSharp.dll
/opt/awsvpnclient/Microsoft.CSharp.dll
/opt/awsvpnclient/Microsoft.VisualBasic.Core.dll
/opt/awsvpnclient/Microsoft.VisualBasic.dll
/opt/awsvpnclient/Microsoft.Win32.Primitives.dll
/opt/awsvpnclient/Microsoft.Win32.Registry.dll
/opt/awsvpnclient/Microsoft.Win32.SystemEvents.dll
/opt/awsvpnclient/Newtonsoft.Json.dll
/opt/awsvpnclient/PangoSharp.dll
/opt/awsvpnclient/Resources
/opt/awsvpnclient/Resources/LINUX-LICENSE.txt
/opt/awsvpnclient/Resources/THIRD-PARTY-LICENSES-GTK.txt
/opt/awsvpnclient/Resources/acvc-64.png
/opt/awsvpnclient/Resources/openvpn
/opt/awsvpnclient/Resources/openvpn/acvc-openvpn
/opt/awsvpnclient/Resources/openvpn/configure-dns
/opt/awsvpnclient/Serilog.Sinks.Console.dll
/opt/awsvpnclient/Serilog.Sinks.File.Header.dll
/opt/awsvpnclient/Serilog.Sinks.File.dll
/opt/awsvpnclient/Serilog.dll
/opt/awsvpnclient/Service
/opt/awsvpnclient/Service/Resources
/opt/awsvpnclient/Service/Resources/openvpn
/opt/awsvpnclient/Service/Resources/openvpn/configure-dns
/opt/awsvpnclient/SymbolicLinkSupport.dll
/opt/awsvpnclient/System.AppContext.dll
/opt/awsvpnclient/System.Buffers.dll
/opt/awsvpnclient/System.CodeDom.dll
/opt/awsvpnclient/System.Collections.Concurrent.dll
/opt/awsvpnclient/System.Collections.Immutable.dll
/opt/awsvpnclient/System.Collections.NonGeneric.dll
/opt/awsvpnclient/System.Collections.Specialized.dll
/opt/awsvpnclient/System.Collections.dll
/opt/awsvpnclient/System.ComponentModel.Annotations.dll
/opt/awsvpnclient/System.ComponentModel.DataAnnotations.dll
/opt/awsvpnclient/System.ComponentModel.EventBasedAsync.dll
/opt/awsvpnclient/System.ComponentModel.Primitives.dll
/opt/awsvpnclient/System.ComponentModel.TypeConverter.dll
/opt/awsvpnclient/System.ComponentModel.dll
/opt/awsvpnclient/System.Configuration.ConfigurationManager.dll
/opt/awsvpnclient/System.Configuration.dll
/opt/awsvpnclient/System.Console.dll
/opt/awsvpnclient/System.Core.dll
/opt/awsvpnclient/System.Data.Common.dll
/opt/awsvpnclient/System.Data.DataSetExtensions.dll
/opt/awsvpnclient/System.Data.SQLite.EF6.dll
/opt/awsvpnclient/System.Data.SQLite.dll
/opt/awsvpnclient/System.Data.SqlClient.dll
/opt/awsvpnclient/System.Data.dll
/opt/awsvpnclient/System.Diagnostics.Contracts.dll
/opt/awsvpnclient/System.Diagnostics.Debug.dll
/opt/awsvpnclient/System.Diagnostics.DiagnosticSource.dll
/opt/awsvpnclient/System.Diagnostics.FileVersionInfo.dll
/opt/awsvpnclient/System.Diagnostics.Process.dll
/opt/awsvpnclient/System.Diagnostics.StackTrace.dll
/opt/awsvpnclient/System.Diagnostics.TextWriterTraceListener.dll
/opt/awsvpnclient/System.Diagnostics.Tools.dll
/opt/awsvpnclient/System.Diagnostics.TraceSource.dll
/opt/awsvpnclient/System.Diagnostics.Tracing.dll
/opt/awsvpnclient/System.Drawing.Common.dll
/opt/awsvpnclient/System.Drawing.Primitives.dll
/opt/awsvpnclient/System.Drawing.dll
/opt/awsvpnclient/System.Dynamic.Runtime.dll
/opt/awsvpnclient/System.Formats.Asn1.dll
/opt/awsvpnclient/System.Formats.Tar.dll
/opt/awsvpnclient/System.Globalization.Calendars.dll
/opt/awsvpnclient/System.Globalization.Extensions.dll
/opt/awsvpnclient/System.Globalization.dll
/opt/awsvpnclient/System.IO.Compression.Brotli.dll
/opt/awsvpnclient/System.IO.Compression.FileSystem.dll
/opt/awsvpnclient/System.IO.Compression.ZipFile.dll
/opt/awsvpnclient/System.IO.Compression.dll
/opt/awsvpnclient/System.IO.FileSystem.AccessControl.dll
/opt/awsvpnclient/System.IO.FileSystem.DriveInfo.dll
/opt/awsvpnclient/System.IO.FileSystem.Primitives.dll
/opt/awsvpnclient/System.IO.FileSystem.Watcher.dll
/opt/awsvpnclient/System.IO.FileSystem.dll
/opt/awsvpnclient/System.IO.IsolatedStorage.dll
/opt/awsvpnclient/System.IO.MemoryMappedFiles.dll
/opt/awsvpnclient/System.IO.Pipes.AccessControl.dll
/opt/awsvpnclient/System.IO.Pipes.dll
/opt/awsvpnclient/System.IO.UnmanagedMemoryStream.dll
/opt/awsvpnclient/System.IO.dll
/opt/awsvpnclient/System.Linq.Expressions.dll
/opt/awsvpnclient/System.Linq.Parallel.dll
/opt/awsvpnclient/System.Linq.Queryable.dll
/opt/awsvpnclient/System.Linq.dll
/opt/awsvpnclient/System.Memory.dll
/opt/awsvpnclient/System.Net.Http.Json.dll
/opt/awsvpnclient/System.Net.Http.dll
/opt/awsvpnclient/System.Net.HttpListener.dll
/opt/awsvpnclient/System.Net.Mail.dll
/opt/awsvpnclient/System.Net.NameResolution.dll
/opt/awsvpnclient/System.Net.NetworkInformation.dll
/opt/awsvpnclient/System.Net.Ping.dll
/opt/awsvpnclient/System.Net.Primitives.dll
/opt/awsvpnclient/System.Net.Quic.dll
/opt/awsvpnclient/System.Net.Requests.dll
/opt/awsvpnclient/System.Net.Security.dll
/opt/awsvpnclient/System.Net.ServicePoint.dll
/opt/awsvpnclient/System.Net.Sockets.dll
/opt/awsvpnclient/System.Net.WebClient.dll
/opt/awsvpnclient/System.Net.WebHeaderCollection.dll
/opt/awsvpnclient/System.Net.WebProxy.dll
/opt/awsvpnclient/System.Net.WebSockets.Client.dll
/opt/awsvpnclient/System.Net.WebSockets.dll
/opt/awsvpnclient/System.Net.dll
/opt/awsvpnclient/System.Numerics.Vectors.dll
/opt/awsvpnclient/System.Numerics.dll
/opt/awsvpnclient/System.ObjectModel.dll
/opt/awsvpnclient/System.Private.CoreLib.dll
/opt/awsvpnclient/System.Private.DataContractSerialization.dll
/opt/awsvpnclient/System.Private.ServiceModel.dll
/opt/awsvpnclient/System.Private.Uri.dll
/opt/awsvpnclient/System.Private.Xml.Linq.dll
/opt/awsvpnclient/System.Private.Xml.dll
/opt/awsvpnclient/System.Reflection.DispatchProxy.dll
/opt/awsvpnclient/System.Reflection.Emit.ILGeneration.dll
/opt/awsvpnclient/System.Reflection.Emit.Lightweight.dll
/opt/awsvpnclient/System.Reflection.Emit.dll
/opt/awsvpnclient/System.Reflection.Extensions.dll
/opt/awsvpnclient/System.Reflection.Metadata.dll
/opt/awsvpnclient/System.Reflection.Primitives.dll
/opt/awsvpnclient/System.Reflection.TypeExtensions.dll
/opt/awsvpnclient/System.Reflection.dll
/opt/awsvpnclient/System.Resources.Reader.dll
/opt/awsvpnclient/System.Resources.ResourceManager.dll
/opt/awsvpnclient/System.Resources.Writer.dll
/opt/awsvpnclient/System.Runtime.CompilerServices.Unsafe.dll
/opt/awsvpnclient/System.Runtime.CompilerServices.VisualC.dll
/opt/awsvpnclient/System.Runtime.Extensions.dll
/opt/awsvpnclient/System.Runtime.Handles.dll
/opt/awsvpnclient/System.Runtime.InteropServices.JavaScript.dll
/opt/awsvpnclient/System.Runtime.InteropServices.RuntimeInformation.dll
/opt/awsvpnclient/System.Runtime.InteropServices.dll
/opt/awsvpnclient/System.Runtime.Intrinsics.dll
/opt/awsvpnclient/System.Runtime.Loader.dll
/opt/awsvpnclient/System.Runtime.Numerics.dll
/opt/awsvpnclient/System.Runtime.Serialization.Formatters.dll
/opt/awsvpnclient/System.Runtime.Serialization.Json.dll
/opt/awsvpnclient/System.Runtime.Serialization.Primitives.dll
/opt/awsvpnclient/System.Runtime.Serialization.Xml.dll
/opt/awsvpnclient/System.Runtime.Serialization.dll
/opt/awsvpnclient/System.Runtime.dll
/opt/awsvpnclient/System.Security.AccessControl.dll
/opt/awsvpnclient/System.Security.Claims.dll
/opt/awsvpnclient/System.Security.Cryptography.Algorithms.dll
/opt/awsvpnclient/System.Security.Cryptography.Cng.dll
/opt/awsvpnclient/System.Security.Cryptography.Csp.dll
/opt/awsvpnclient/System.Security.Cryptography.Encoding.dll
/opt/awsvpnclient/System.Security.Cryptography.OpenSsl.dll
/opt/awsvpnclient/System.Security.Cryptography.Pkcs.dll
/opt/awsvpnclient/System.Security.Cryptography.Primitives.dll
/opt/awsvpnclient/System.Security.Cryptography.ProtectedData.dll
/opt/awsvpnclient/System.Security.Cryptography.X509Certificates.dll
/opt/awsvpnclient/System.Security.Cryptography.Xml.dll
/opt/awsvpnclient/System.Security.Cryptography.dll
/opt/awsvpnclient/System.Security.Permissions.dll
/opt/awsvpnclient/System.Security.Principal.Windows.dll
/opt/awsvpnclient/System.Security.Principal.dll
/opt/awsvpnclient/System.Security.SecureString.dll
/opt/awsvpnclient/System.Security.dll
/opt/awsvpnclient/System.ServiceModel.Http.dll
/opt/awsvpnclient/System.ServiceModel.Primitives.dll
/opt/awsvpnclient/System.ServiceModel.Web.dll
/opt/awsvpnclient/System.ServiceModel.dll
/opt/awsvpnclient/System.ServiceProcess.dll
/opt/awsvpnclient/System.Text.Encoding.CodePages.dll
/opt/awsvpnclient/System.Text.Encoding.Extensions.dll
/opt/awsvpnclient/System.Text.Encoding.dll
/opt/awsvpnclient/System.Text.Encodings.Web.dll
/opt/awsvpnclient/System.Text.Json.dll
/opt/awsvpnclient/System.Text.RegularExpressions.dll
/opt/awsvpnclient/System.Threading.Channels.dll
/opt/awsvpnclient/System.Threading.Overlapped.dll
/opt/awsvpnclient/System.Threading.Tasks.Dataflow.dll
/opt/awsvpnclient/System.Threading.Tasks.Extensions.dll
/opt/awsvpnclient/System.Threading.Tasks.Parallel.dll
/opt/awsvpnclient/System.Threading.Tasks.dll
/opt/awsvpnclient/System.Threading.Thread.dll
/opt/awsvpnclient/System.Threading.ThreadPool.dll
/opt/awsvpnclient/System.Threading.Timer.dll
/opt/awsvpnclient/System.Threading.dll
/opt/awsvpnclient/System.Transactions.Local.dll
/opt/awsvpnclient/System.Transactions.dll
/opt/awsvpnclient/System.ValueTuple.dll
/opt/awsvpnclient/System.Web.HttpUtility.dll
/opt/awsvpnclient/System.Web.dll
/opt/awsvpnclient/System.Windows.Extensions.dll
/opt/awsvpnclient/System.Windows.dll
/opt/awsvpnclient/System.Xml.Linq.dll
/opt/awsvpnclient/System.Xml.ReaderWriter.dll
/opt/awsvpnclient/System.Xml.Serialization.dll
/opt/awsvpnclient/System.Xml.XDocument.dll
/opt/awsvpnclient/System.Xml.XPath.XDocument.dll
/opt/awsvpnclient/System.Xml.XPath.dll
/opt/awsvpnclient/System.Xml.XmlDocument.dll
/opt/awsvpnclient/System.Xml.XmlSerializer.dll
/opt/awsvpnclient/System.Xml.dll
/opt/awsvpnclient/System.dll
/opt/awsvpnclient/Tmds.DBus.dll
/opt/awsvpnclient/WindowsBase.dll
/opt/awsvpnclient/de
/opt/awsvpnclient/de/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/es
/opt/awsvpnclient/es/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/fr
/opt/awsvpnclient/fr/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/it
/opt/awsvpnclient/it/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/ja
/opt/awsvpnclient/ja/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/ko
/opt/awsvpnclient/ko/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/libSystem.Globalization.Native.so
/opt/awsvpnclient/libSystem.IO.Compression.Native.so
/opt/awsvpnclient/libSystem.Native.so
/opt/awsvpnclient/libSystem.Net.Security.Native.so
/opt/awsvpnclient/libSystem.Security.Cryptography.Native.OpenSsl.so
/opt/awsvpnclient/libclrgc.so
/opt/awsvpnclient/libclrjit.so
/opt/awsvpnclient/libcoreclr.so
/opt/awsvpnclient/libhostfxr.so
/opt/awsvpnclient/libhostpolicy.so
/opt/awsvpnclient/mscorlib.dll
/opt/awsvpnclient/netstandard.dll
/opt/awsvpnclient/pt-BR
/opt/awsvpnclient/pt-BR/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/zh-Hans
/opt/awsvpnclient/zh-Hans/AWSVPNClient.Core.resources.dll
/opt/awsvpnclient/zh-Hant
/opt/awsvpnclient/zh-Hant/AWSVPNClient.Core.resources.dll
/usr/lib/systemd/system-preset/70-awsvpnclient.preset
/usr/lib/systemd/system/awsvpnclient.service
/usr/share/applications/awsvpnclient.desktop
/usr/share/doc/awsvpnclient
/usr/share/doc/awsvpnclient/copyright
/usr/share/pixmaps/acvc-64.png
